# relay/trading-relay.service
[Unit]
Description=Trading Assistant Relay
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/trading-relay
ExecStart=/opt/trading-relay/start.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
